# flutter-env.service
[Unit]
Description=Flutter Development Environment
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/start-flutter-env.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target